# src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisord
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

;; ==========================================================
;; Gitaly — Git RPC backend (must start first)
;; ==========================================================
[program:gitaly]
command=/home/git/gitaly/_build/bin/gitaly /home/git/gitaly/config.toml
user=git
directory=/home/git/gitaly
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisord/gitaly-stdout.log
stderr_logfile=/var/log/supervisord/gitaly-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/git",PATH="/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"

;; ==========================================================
;; GitLab Workhorse — Smart reverse proxy
;; ==========================================================
[program:workhorse]
command=/home/git/gitlab-workhorse/gitlab-workhorse
    -listenUmask 0
    -listenNetwork unix
    -listenAddr /home/git/gitlab/tmp/sockets/gitlab-workhorse.socket
    -authBackend http://127.0.0.1:8080
    -authSocket /home/git/gitlab/tmp/sockets/gitlab.socket
    -documentRoot /home/git/gitlab/public
    -secretPath /home/git/gitlab/.gitlab_workhorse_secret
user=git
directory=/home/git/gitlab
autostart=true
autorestart=true
priority=150
stdout_logfile=/var/log/supervisord/workhorse-stdout.log
stderr_logfile=/var/log/supervisord/workhorse-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/git",PATH="/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"

;; ==========================================================
;; Puma — Rails application server
;; ==========================================================
[program:puma]
command=bundle exec puma -C /home/git/gitlab/config/puma.rb
user=git
directory=/home/git/gitlab
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisord/puma-stdout.log
stderr_logfile=/var/log/supervisord/puma-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/git",RAILS_ENV="production",NODE_ENV="production",PATH="/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin"

;; ==========================================================
;; Sidekiq — Background job processor
;; ==========================================================
[program:sidekiq]
command=bundle exec sidekiq -C /home/git/gitlab/config/sidekiq_queues.yml -e production
user=git
directory=/home/git/gitlab
autostart=true
autorestart=true
priority=300
stdout_logfile=/var/log/supervisord/sidekiq-stdout.log
stderr_logfile=/var/log/supervisord/sidekiq-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=HOME="/home/git",RAILS_ENV="production",NODE_ENV="production",PATH="/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin",MALLOC_ARENA_MAX="2"

;; ==========================================================
;; Nginx — Frontend HTTP proxy
;; ==========================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=400
stdout_logfile=/var/log/supervisord/nginx-stdout.log
stderr_logfile=/var/log/supervisord/nginx-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

;; ==========================================================
;; SSHD — Git over SSH
;; ==========================================================
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
priority=500
stdout_logfile=/var/log/supervisord/sshd-stdout.log
stderr_logfile=/var/log/supervisord/sshd-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
